# tests/dynalite/supervisord.conf
[unix_http_server]
file = /tmp/dynamodb_supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/dynamodb_supervisor.sock

[supervisord]

[program:dynamodb_local]
command=java -Djava.library.path=./DynamoDBLocal_lib -jar DynamoDBLocal.jar -inMemory
directory=dynamodb_local
autostart=false

[program:dynalite]
command=dynalite
autostart=false